# cron.hourly/depscan
#!/bin/sh

# force Rebuild Gentoo init.d service dependency tree.
depscan.sh --update
