# local.d/01-swap-esc-caps.start
#!/bin/sh

# swap Esc and CapsLock, for vimmer !
echo "keycode 58 = Escape" | loadkeys -
