# cron.daily/localepurge
#!/bin/sh

/usr/sbin/localepurge.config

localepurge
