# cron.daily/portage
#! /bin/sh

# eix-sync = "emerge --sync" + "eix-update"
eix-sync

